# supervisord.conf
[program:yak-claude-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work database --queue=yak-claude --timeout=600 --tries=3 --sleep=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=storage/logs/yak-claude-worker.log
stopwaitsecs=610

[program:default-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work database --queue=default --timeout=30 --tries=3 --sleep=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=storage/logs/default-worker.log
stopwaitsecs=40
